# ansible/roles/oracle_catcher/templates/oracle_catcher.service.j2
[Unit]
Description=Oracle ARM Capacity Catcher
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
Environment="PYTHONUNBUFFERED=1"
WorkingDirectory=/opt/oracle_catcher
ExecStart=/opt/oracle_catcher/venv/bin/python3 /opt/oracle_catcher/oracle_turbo_catcher.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Email для уведомлений (из Ansible Vault через EnvironmentFile)
EnvironmentFile=-/opt/oracle_catcher/.env

[Install]
WantedBy=multi-user.target
